# fabsettings/templates/gunicorn.conf
[program:%(project)s]
command=%(supervisor_dir)s/gunicorn_start_%(project)s.bash
directory=%(project_dir)s
user=%(user)s
autostart=true
autorestart=true
redirect_stderr=True
